# initrd/emergency.service.d/50-sulogin-force.conf
# Tell sulogin to start a root shell if password data is not available.
[Service]
ExecStart=
ExecStart=-/bin/sh -c "/sbin/sulogin --force; /usr/bin/systemctl --fail --no-block default"
